# deploy/systemd/press-ingest.timer
[Unit]
Description=Run press ingest five times daily

[Timer]
Timezone=Asia/Seoul
OnCalendar=*-*-* 09:00:00
OnCalendar=*-*-* 13:00:00
OnCalendar=*-*-* 16:00:00
OnCalendar=*-*-* 18:00:00
OnCalendar=*-*-* 23:00:00
Persistent=true
Unit=press-ingest.service

[Install]
WantedBy=timers.target
