# deploy/systemd/press-ingest.timer
[Unit]
Description=Run press ingest five times daily

[Timer]
Timezone=Asia/Seoul
OnCalendar=*-*-* 09:30:00
OnCalendar=*-*-* 13:30:00
OnCalendar=*-*-* 15:30:00
OnCalendar=*-*-* 17:30:00
OnCalendar=*-*-* 00:00:00
Persistent=true
Unit=press-ingest.service

[Install]
WantedBy=timers.target
