# deploy/systemd/press-dashboard.service
[Unit]
Description=Press Dashboard HTTP Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/wjdgns429/press-dashboard
ExecStart=/home/wjdgns429/press-dashboard/.venv/bin/python local_dashboard.py --db-path press_unified.db --host 0.0.0.0 --port 80
Restart=always
RestartSec=5
User=wjdgns429
Group=wjdgns429
Environment=PYTHONUNBUFFERED=1
StandardOutput=append:/home/wjdgns429/press-dashboard/logs/dashboard-systemd.log
StandardError=append:/home/wjdgns429/press-dashboard/logs/dashboard-systemd.log

[Install]
WantedBy=multi-user.target
